# ports/sysutils/net.macports.buildworker/files/net.macports.buildworker.base.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>EnvironmentVariables</key>
        <dict>
            <key>HOME</key>
            <string>@PREFIX@/var/buildworker/base</string>
        </dict>
        <key>KeepAlive</key>
        <dict>
            <key>SuccessfulExit</key>
            <false/>
        </dict>
        <key>Label</key>
        <string>net.macports.buildworker.base</string>
        <key>ProgramArguments</key>
        <array>
            <string>@PREFIX@/bin/twistd-2.7</string>
            <string>--nodaemon</string>
            <string>--no_save</string>
            <string>--python=buildbot.tac</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>UserName</key>
        <string>buildbot</string>
        <key>WorkingDirectory</key>
        <string>@PREFIX@/var/buildworker/base</string>
    </dict>
</plist>
